# docker/php/supervisord.conf
[supervisord]
nodaemon=true
user=app
loglevel=debug

[program:php-fpm]
command=php-fpm